# src-tauri/extra/bookos-battery-adaptive.timer
[Unit]
Description=BookOS Adaptive Charging timer — runs every minute
Requires=bookos-battery-adaptive.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
