# src/BookOS-Settings-0.4.0/src-tauri/extra/bookos-battery-logger.service
[Unit]
Description=BookOS Battery Logger — records one CSV sample
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/bookos/bookos-battery-logger.sh
