# src/BookOS-Settings-0.4.0/src-tauri/extra/search/bookos-search-watcher.service
[Unit]
Description=BookOS Semantic Search — filesystem watcher (indexado incremental)
After=default.target

[Service]
Type=simple
Environment=SENTENCE_TRANSFORMERS_HOME=/opt/bookos-search/models
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=/bin/bash /opt/bookos-search/watcher.sh
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=default.target
