# src/BookOS-Settings-0.4.0/src-tauri/extra/search/bookos-search-reindex.service
[Unit]
Description=BookOS Semantic Search — reindex total

[Service]
Type=oneshot
Environment=SENTENCE_TRANSFORMERS_HOME=/opt/bookos-search/models
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=/opt/bookos-search/venv/bin/python /opt/bookos-search/indexador.py --full
Nice=15
IOSchedulingClass=idle
