# src/BookOS-Settings-0.4.0/src-tauri/extra/search/bookos-search-reindex.timer
[Unit]
Description=BookOS Semantic Search — reindex semanal

[Timer]
OnBootSec=10min
OnUnitActiveSec=1w
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
